# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0

[program:couchdb]
command=/opt/couchdb/bin/couchdb
user=couchdb
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=COUCHDB_USER="%(ENV_COUCHDB_USER)s",COUCHDB_PASSWORD="%(ENV_COUCHDB_PASSWORD)s"

[program:mcp]
command=node /app/dist/main.js
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=COUCHDB_URL="http://localhost:5984",COUCHDB_USER="%(ENV_COUCHDB_USER)s",COUCHDB_PASSWORD="%(ENV_COUCHDB_PASSWORD)s",COUCHDB_DATABASE="%(ENV_COUCHDB_DATABASE)s",VAULT_NAME="%(ENV_VAULT_NAME)s",PORT="8787"
startsecs=5
